# server/cron/crontab
0 0 * * * bash /home/centos/backup.sh
0 0 * * 0 bash /home/centos/cleanup.sh